# cfport.service
[Unit]
Description=cfporttest service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/app
Environment=LISTEN_ADDR=:8888
ExecStart=/app/cfporttest
Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
